# deploy/hrrr-server.service
[Unit]
Description=HRRR forecast cache HTTP server
After=network.target

[Service]
Type=simple
User=kiteagent
Group=kiteagent
WorkingDirectory=/opt/kiteagent
ExecStart=/opt/kiteagent/hrrr-server /opt/kiteagent/hrrr.db
Restart=always
RestartSec=5
Environment=RUST_LOG=info
Environment=BIND=0.0.0.0:8081

[Install]
WantedBy=multi-user.target
