# deploy/kiteagent-agent.service
[Unit]
Description=KiteAgent Weather Daemon
After=network.target kiteagent-server.service
Requires=kiteagent-server.service

[Service]
Type=simple
User=kiteagent
Group=kiteagent
WorkingDirectory=/app/kiteagent
ExecStart=/app/kiteagent/target/release/kiteagent-agent /app/kiteagent/config.toml
Restart=always
RestartSec=5
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
