# deploy/kiteagent-server.service
[Unit]
Description=KiteAgent Push Server
After=network.target hrrr-server.service live-server.service

[Service]
Type=simple
User=kiteagent
Group=kiteagent
WorkingDirectory=/app/kiteagent
ExecStart=/app/kiteagent/target/release/kiteagent-server
Restart=always
RestartSec=5
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
